# systemd/captive-portal.service
[Unit]
Description=VPN-AP Captive Portal Web Interface
After=network.target hostapd.service dnsmasq.service
Wants=hostapd.service dnsmasq.service

[Service]
Type=simple

# Load configuration including PROJECT_DIR
EnvironmentFile=-/etc/default/vpn-ap

# Use PROJECT_DIR to locate script (defaults to /home/pi/VPN-AP if not set)
ExecStart=/bin/bash -c '/usr/bin/python3 ${PROJECT_DIR:-/home/pi/VPN-AP}/scripts/captive-portal-server.py'
Restart=on-failure
RestartSec=5
User=root

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
