# systemd/captive-portal.service
[Unit]
Description=VPN-AP Captive Portal Web Interface
After=network.target hostapd.service dnsmasq.service
Wants=hostapd.service dnsmasq.service

[Service]
Type=simple

# Load configuration including PROJECT_DIR
EnvironmentFile=-/etc/default/vpn-ap

# PROJECT_DIR is set by /etc/default/vpn-ap (created during setup)
ExecStart=/bin/bash -c '/usr/bin/python3 ${PROJECT_DIR:?PROJECT_DIR not set - run setup.sh}/scripts/captive-portal-server.py'
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=10
User=root

# Hardening: restrict what root can do even if the portal is compromised
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/var/lib/vpn-ap /var/log /etc/dnsmasq.d /run
NoNewPrivileges=no
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
