# systemd/vpn-ap-watchdog.service
[Unit]
Description=VPN-AP Watchdog Service
After=network.target

[Service]
Type=oneshot

# Default values (overridden by EnvironmentFile)
Environment=VPN_INTERFACE=wg0
Environment=UPSTREAM_INTERFACES="eth0 wlan0"

# Load configuration including PROJECT_DIR
EnvironmentFile=-/etc/default/vpn-ap

# Use installed script path (setup.sh copies to /usr/local/bin/)
ExecStart=/usr/local/bin/vpn-ap-watchdog
TimeoutStartSec=120
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
